# .config/systemd/user/syncthing.service
[Unit]
Description=SyncThing service

[Service]
ExecStart=/usr/bin/syncthing

[Install]
WantedBy=default.target
