# .config/systemd/user/vdirsyncer.service
[Unit]
Description=Calendar and contact synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
